# roles/postgresql/templates/postgresql.service.j2
.include /lib/systemd/system/postgresql-{{postgresql_version}}.service
[Service]
Environment=PGDATA={{ postgresql_conf_directory }}